# backend/launchd/com.stocks_scan.warm.us.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  美股盘后缓存刷新 + 双扫描.
  运行时间: 每日本地时间 22:30 (≈ 美股收盘 16:00 ET +30 分钟, 假设系统 TZ=CEST/CET)
  夏令 CEST(+2) 对应 EDT(-4) 收盘 22:00 本地;
  冬令 CET(+1) 对应 EST(-5) 收盘 22:00 本地; 两种情况 22:30 皆可覆盖.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stocks_scan.warm.us</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/robin/Documents/Stocks_Scan/backend/.venv/bin/python</string>
        <string>/Users/robin/Documents/Stocks_Scan/backend/warm.py</string>
        <string>--market</string>
        <string>us</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/robin/Documents/Stocks_Scan/backend</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/robin/Documents/Stocks_Scan/backend/logs/warm.us.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/robin/Documents/Stocks_Scan/backend/logs/warm.us.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>STOCKS_SCAN_CACHE</key>
        <string>/Users/robin/Documents/Stocks_Scan/backend/cache</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
